# com.nubhq.video-processor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nubhq.video-processor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mattcolvin/nubhq/api/venv/bin/python</string>
        <string>-m</string>
        <string>app.worker.intelligent_processor</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mattcolvin/nubhq/api</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NUBHQ_INPUT</key>
        <string>/Volumes/Lil Hoss/NubHQ/01_Inbox_ToSort</string>
        <key>NUBHQ_OUTPUT</key>
        <string>/Volumes/Big Hoss/NubHQ/03_Exports</string>
        <key>NUBHQ_LIBRARY</key>
        <string>/Volumes/Big Hoss/NubHQ/02_Library</string>
        <key>NUBHQ_ARCHIVE</key>
        <string>/Volumes/Big Hoss/NubHQ/04_Archive</string>
        <key>NUBHQ_DATA</key>
        <string>/Volumes/Big Hoss/NubHQ/database</string>
        <key>NUBHQ_TEMP</key>
        <string>/Volumes/Lil Hoss/NubHQ/temp</string>
        <key>NUBHQ_QUICK_MODE</key>
        <string>youtube</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>PathState</key>
        <dict>
            <key>/Volumes/Lil Hoss</key>
            <true/>
        </dict>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Volumes/Lil Hoss/NubHQ/logs/processor.log</string>
    <key>StandardErrorPath</key>
    <string>/Volumes/Lil Hoss/NubHQ/logs/processor-error.log</string>
</dict>
</plist>
